# docker/supervisord.conf
# docker/supervisord.conf
# Supervisor Configuration für HCC Plan Multi-User Sessions
# Verwaltet X11 Display Server, VNC Server und HCC Plan Application

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_RESOLUTION)s -ac +extension GLX +render -noreset
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log
priority=10

[program:x11vnc] 
command=x11vnc -display %(ENV_DISPLAY)s -forever -usepw -shared -rfbport %(ENV_VNC_PORT)s -rfbauth /root/.vnc/passwd
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log
priority=20
depends_on=xvfb

[program:xfce4-session]
command=xfce4-session
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xfce4.log
stderr_logfile=/var/log/supervisor/xfce4.log
environment=DISPLAY="%(ENV_DISPLAY)s"
priority=30
depends_on=xvfb

[program:hcc-plan-app]
command=python main.py --project-filter=%(ENV_HCC_PROJECT_ID)s --user-id=%(ENV_HCC_USER_ID)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/hcc-plan.log
stderr_logfile=/var/log/supervisor/hcc-plan.log
environment=DISPLAY="%(ENV_DISPLAY)s",PYTHONPATH="/app",DATABASE_URL="%(ENV_DATABASE_URL)s"
priority=40
depends_on=xfce4-session
user=root